# systemd/isidore-tmux.service
[Unit]
Description=Isidore tmux session (persistent interactive Claude)
After=network.target

[Service]
Type=forking
User=isidore
Group=isidore
ExecStart=/usr/bin/tmux new-session -d -s isidore
ExecStop=/usr/bin/tmux kill-session -t isidore
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
